# scripts/uwsgi.conf
description "uWSGI server for Graco"

start on started mysql 
stop on starting rc RUNLEVEL=[016]

respawn

exec /usr/local/bin/uwsgi -c /home/graco/pump/WebPilot/uwsgi-graco.ini
